# threat-intel/tests/clean_files/systemd.service
[Unit]
Description=Example Service
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/example --config /etc/example.conf
Restart=on-failure

[Install]
WantedBy=multi-user.target
